# config/systemd/.config/systemd/user/check-battery.service
[Unit]
Description=Check if battery is low

[Service]
ExecStart=/home/cat/.bin/check-battery.sh
