# config/systemd/.config/systemd/user/check-battery.timer
[Unit]
Description=Timer to check if battery level is low

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
